# systemd/gauls-copy-trader.service
[Unit]
Description=Gauls Main Copy Trader
Documentation=https://github.com/gauls-trading-system
After=network.target gauls-telegram-listener.service
PartOf=gauls-trading.target
Requires=gauls-trading.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/gauls-copy-trading-system
Environment="PATH=/gauls-copy-trading-system/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/gauls-copy-trading-system"
ExecStart=/gauls-copy-trading-system/venv/bin/python /gauls-copy-trading-system/gauls_copy_trader.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=5

# Resource limits
MemoryLimit=1G
CPUQuota=75%

# Logging
StandardOutput=append:/gauls-copy-trading-system/logs/copy_trader.log
StandardError=append:/gauls-copy-trading-system/logs/copy_trader_error.log

# Health check
ExecStartPre=/bin/bash -c 'mkdir -p /gauls-copy-trading-system/logs'

[Install]
WantedBy=gauls-trading.target